# ansible/roles/deadeye/files/deadeye-daemon.service
[Unit]
Description=Deadeye Vision System Daemon

[Service]
ExecStart=/usr/local/bin/deadeyed

[Install]
WantedBy=multi-user.target
