# ansible/roles/deadeye/files/mnt-deadeye.mount
[Unit]
Description=Deadeye Log Storage

[Mount]
What=/dev/sda1
Where=/mnt/deadeye
