# ansible/roles/deadeye/files/var-opt-deadeye.mount
[Unit]
Description=Deadeye Log Storage

[Mount]
What=/dev/sda1
Where=/var/opt/deadeye
