# ansible/roles/deadeye/templates/var-opt-deadeye.mount.j2
[Unit]
Description=Deadeye Log Storage

[Mount]
What=/dev/sda1
Where={{ deadeye_mount }}
